# scripts/restash-refresh.timer
[Unit]
Description=Run Restash Quick Refresh daily

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
